# conf/solr.supervisor.conf
[program:solr]
directory = /opt/solr/example
command = java -jar start.jar  -Djava.net.preferIPv4Stack=true -Djava.net.preferIPv4Addresses=true -OPTIONS=jsp
autostart = true
autorestart = true

